# app/templates/app_supervisor.ini
[unix_http_server]
file={{config.user_home}}/var/supervisor.sock


[supervisord]
directory={{config.user_home}}
pidfile={{config.user_home}}/var/run/supervisord.pid
logfile_backups=1
logfile = {{config.user_home}}/logs/supervisord.log
logfile_maxbytes = {{supervisord.logfile_maxbytes}}
loglevel = {{supervisord.loglevel}}
environment=
identifier=supervisor


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl=unix://{{config.user_home}}/var/supervisor.sock


[program:uwsgi]
pidfile={{config.user_home}}/var/run/uwsgi.pid
command = uwsgi --ini %(here)s/uwsgi.ini

stdout_logfile={{config.user_home}}/logs/uwsgi.log
stderr_logfile={{config.user_home}}/logs/uwsgi.err
redirect_stderr=true

autostart=true
autorestart=true

startsecs=10
stopsignal=QUIT


[program:nginx]
pidfile={{config.user_home}}/var/run/nginx.pid
command = nginx -p ~ -c %(here)s/nginx.conf

stdout_logfile={{config.user_home}}/logs/nginx.log
stderr_logfile={{config.user_home}}/logs/nginx.err
redirect_stderr=true

autostart=true
autorestart=true

stopsignal=INT
startsecs=10
