# app/templates/rc.local
#!/bin/sh -e

python /etc/init_supervisors.py

exit 0
